# lampserver/lampserver.socket
[Socket]
ListenStream=8080
KeepAlive=true

[Install]
WantedBy=sockets.target
